# systemd/transmission_banner.service
[Unit]
Description=Run Transmission to update issue banner

[Service]
Type=oneshot
ExecStart=/usr/bin/transmission update-banner
