# systemd/transmission.service
[Unit]
Description=Run Transmission

[Service]
Type=oneshot
ExecStart=/usr/bin/transmission
